# scripts/.local/bin/cron.hourly/hourly-chime
#!/bin/sh
exec dunstify -t 0 -r 3454 "Hour Chime: $(date +%I:%M\ %p)" "Another hour has gone by"
